# systemd/twix-rsync.service
[Unit]
Description=Rsync disk.img contents to /twixfiles

[Service]
Type=oneshot
ExecStart=/usr/local/bin/twix-rsync
User=root
Group=root
