# systemd/usb-rsync.service
[Unit]
Description=Rsync disk.img contents to /usbfiles

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-rsync
User=root
Group=root
